# ceph/templates/default/upstart-ceph-osd-all.conf.erb
description "Ceph OSD (start all instances)"

start on filesystem

task

script
  find /srv/ -mindepth 1 -maxdepth 1 -regextype posix-egrep -regex '.*/osd\.[0-9]+' -printf '%P\n' \
  | while read f; do
    if [ -e "/srv/$f/done" ]; then
        id="${f##osd.}"

	# upstart start(8) fails if the job is already running
	# https://bugs.launchpad.net/upstart/+bug/878322
	status ceph-osd OSD_ID="$id" \
	| {
	  set -e
	  read name goal_status rest
	  case "$goal_status" in
	    start/*)
	      # no action needed
	      ;;
	    *)
	      start ceph-osd OSD_ID="$id"
	      ;;
	  esac
	}
    fi
  done
end script
